# JETSON_ROBOT/systemd/strew-robot-agent.service
[Unit]
Description=STREW Jetson Robot Agent
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/jetson/STREW_VISION/strew-jetson
EnvironmentFile=/home/jetson/STREW_VISION/strew-jetson/.env
ExecStart=/home/jetson/STREW_VISION/strew-jetson/.venv/bin/python -m main
Restart=always
RestartSec=3
User=jetson

[Install]
WantedBy=multi-user.target
